# examples/aws/poc/dsf_deployment_for_saas/services_scripts/send_health.service
[Unit]
Description=Send health check ro SaaS

[Service]
Type=oneshot
ExecStart=/usr/bin/stdbuf -oL /usr/local/bin/send_health.sh
Environment="JSONAR_LOCALDIR=/imperva/local"